# docker/local/services/supervisor/conf.d/schedule.conf
[program:laravel-schedule]
command=/usr/local/bin/php /var/www/api/artisan schedule:work
autostart=true
autorestart=true